# templates/service_clamsmtpd.j2
#!/bin/ksh

daemon="/usr/local/sbin/clamsmtpd"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1

